# systemd/proxy-to-livestatus.socket
[Unit]
StopWhenUnneeded=true

[Socket]
ListenStream=6557
